# pkgbuild/biglinux-tlp-config.install
## arg 1:  the new package version
#pre_install() {
	# do something here
#}

## arg 1:  the old package version
post_remove() {
	systemctl disable --now tlp.service
}
